# whisper.service
[Unit]
Description=Lab-Endoc Whisper Transcription Service (Docker)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=dev
Group=dev
WorkingDirectory=/home/dev/labendoc/transcribe
RemainAfterExit=yes

ExecStart=/usr/bin/docker-compose -f /home/dev/labendoc/transcribe/docker-compose.yml up -d
ExecStop=/usr/bin/docker-compose -f /home/dev/labendoc/transcribe/docker-compose.yml down
ExecReload=/usr/bin/docker-compose -f /home/dev/labendoc/transcribe/docker-compose.yml restart

StandardOutput=journal
StandardError=journal
SyslogIdentifier=whisper-docker

[Install]
WantedBy=multi-user.target 